# roles/mmdvmhost/files/mmdvmhost@.service
[Unit]
Description=Multi Mode Digital Voice Modem Host for %i

[Service]
ExecStart=/usr/bin/MMDVMHost /etc/mmdvmhost-%i.ini

[Install]
WantedBy=multi-user.target
