# roles/git_mirror/files/git-mirror.service
[Unit]
Description=Update Git mirrors

[Service]
Type=oneshot
User=git-mirror
Group=git-mirror
ExecStart=/usr/bin/git-mirror
